# infra/supervisord.conf
[supervisord]
# Stay foreground
nodaemon=true

[program:influxdb]
command=/usr/bin/influxdb -pidfile /opt/influxdb/shared/influxdb.pid -config /opt/influxdb/shared/config.toml
